# buildroot/board/tcrpfriend/rootfs-overlay/etc/udev/rules.d/99-custom.rules
KERNEL=="sda", ACTION=="add", RUN+="/bin/bash /root/load-sda-first.sh"
