# buildroot/board/tcrpfriend/rootfs-overlay/etc/udev/rules.d/99-custom.rules
KERNEL=="loop: module loaded", ACTION=="add", RUN+="/bin/sh /root/load-sda-first.sh"
